# startups/webserver.ini
[program:wsgi]
directory=%(here)s/wsgi-app
command=gunicorn app:app --workers=2 --bind=0.0.0.0:9001
stdout_logfile=%(here)s/logs/wsgi-out.log
stderr_logfile=%(here)s/logs/wsgi-err.log

[program:flask]
directory=%(here)s/flask-app
command=gunicorn app:app --workers=2 --bind=0.0.0.0:9002
stdout_logfile=%(here)s/logs/flask-out.log
stderr_logfile=%(here)s/logs/flask-err.log

[program:django]
directory=%(here)s/django-app
command=gunicorn app:app --workers=2 --bind=0.0.0.0:9003
stdout_logfile=%(here)s/logs/django-out.log
stderr_logfile=%(here)s/logs/django-err.log

[program:falcon]
directory=%(here)s/falcon-app
command=gunicorn app:app --workers=2 --bind=0.0.0.0:9004
stdout_logfile=%(here)s/logs/falcon-out.log
stderr_logfile=%(here)s/logs/falcon-err.log
